# config/systemd/niello-networking.service
[Unit]
Description=Niello Networking & Bluetooth Setup
After=network-pre.target
Before=network-online.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/bash -c '
set -euo pipefail

# Ensure NetworkManager is enabled and started
systemctl enable NetworkManager 2>/dev/null || true
systemctl start NetworkManager 2>/dev/null || true

# Wait for NetworkManager to report full online state (not just started)
# This brings up network-online.target for us
nm-online -s --timeout=30 2>/dev/null || {
    echo "niello-networking: nm-online timed out, continuing anyway"
}

# Ensure bluetooth service is enabled and started
systemctl enable bluetooth 2>/dev/null || true
systemctl start bluetooth 2>/dev/null || true
'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
