# config/systemd/niello-boot-upgrade.service
[Unit]
Description=Niello bootc upgrade check before session start
After=niello-networking.service
Before=greetd.service
DefaultDependencies=no
ConditionPathExists=!/etc/niello-skip-boot-upgrade

[Service]
Type=oneshot
ExecStart=/usr/local/bin/niello-boot-upgrade.sh
TimeoutStartSec=30
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
