# config/systemd/ollama.service
[Unit]
Description=Ollama local LLM server
After=network.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ollama serve
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical-session.target